# src/test/load/hub-locust.conf
description "hub locust tester"
# This is version 0.1.5a!

start on runlevel [2345]
stop on shutdown


env DEPLOYED_DIR=/home/ubuntu/load
env APP_PATH=/home/ubuntu/load
env APP_NAME=hub-locust


pre-start script
if [ -f  /home/ubuntu/load/prestart.sh ]
then
exec sudo -u ubuntu -H  /home/ubuntu/load/prestart.sh
fi
end script

post-start script
if [ -f  /home/ubuntu/load/poststart.sh ]
then
exec sudo -u ubuntu -H  /home/ubuntu/load/poststart.sh
fi
end script

#
# Upstart attempts to track *the first pid in the exec command*
# (unless overriden by "expect" commands.
# In our case, just run the script *and never background anything*
# i.e. Just run commands, don't attempt shell process control with nohup or &
#
script
#exec 2>>/dev/.initramfs/hub-locust.log
#set -x
exec sudo -u ubuntu -H /home/ubuntu/load/start.sh
end script

pre-stop script
if [ -f  /home/ubuntu/load/prestop.sh ]
then
exec sudo -u ubuntu -H  /home/ubuntu/load/prestop.sh
fi
end script

post-stop script
if [ -f  /home/ubuntu/load/poststop.sh ]
then
exec sudo -u ubuntu -H  /home/ubuntu/load/poststop.sh
fi
end script